# leoxgpon-scraper.service
[Unit]
Description=LeoX GPON metrics scraper
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/leoxgpon
ExecStart=/usr/bin/python3 /opt/leoxgpon/scraper.py \
    --interval 60 \
    --port 9101 \
    --db /var/lib/leoxgpon/leoxgpon.db \
    --no-files
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
